# config/com.user.last-timestamp-reader.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.user.last-timestamp-reader</string>
    <key>ProgramArguments</key>
    <array>
        <string>/Users/alienspirit/Documents/25D/L2/Repo4/venv/bin/python3</string>
        <string>/Users/alienspirit/Documents/25D/L2/Repo4/bin/last-timestamp-reader.py</string>
    </array>
    <key>WorkingDirectory</key>
    <string>/Users/alienspirit/Documents/25D/L2/Repo4</string>
    <key>RunAtLoad</key>
    <true/>
    <key>KeepAlive</key>
    <true/>
    <key>StandardOutPath</key>
    <string>/Users/alienspirit/Documents/25D/L2/Repo4/logs/last-timestamp-reader.log</string>
    <key>StandardErrorPath</key>
    <string>/Users/alienspirit/Documents/25D/L2/Repo4/logs/last-timestamp-reader-error.log</string>
    <key>EnvironmentVariables</key>
    <dict>
        <key>PATH</key>
        <string>/Users/alienspirit/Documents/25D/L2/Repo4/venv/bin:/usr/bin:/bin</string>
    </dict>
</dict>
</plist>
